# deploy/flask_rest_api.ini
[program:flask_rest_api]
chdir=/home/bewi/flask_rest_api
home=/home/bewi/flask_rest_api
pythonpath=/home/bewi/flask_rest_api/ENV/bin/python
virtualenv=/home/bewi/flask_rest_api/ENV+
command=/bin/bash -c '~/flask_rest_api/ENV/bin/uwsgi uwsgi.ini'
directory=/home/bewi/flask_rest_api
autostart=yes
autorestart=yes
stopsignal=INT
stopasgroup=true
killasgroup=true

