# deploy/flask_api.ini
[program:flask_api]
chdir=/home/bewi/flask_api
home=/home/bewi/flask_api
pythonpath=/home/bewi/flask_api/ENV/bin/python
virtualenv=/home/bewi/flask_api/ENV+
command=/bin/bash -c '~/flask_api/ENV/bin/uwsgi uwsgi.ini'
directory=/home/bewi/flask_api
autostart=yes
autorestart=yes
stopsignal=INT
stopasgroup=true
killasgroup=true

